# packages/myos-settings/myos-ai-tools-setup.service
[Unit]
Description=Verify MyOS bundled AI runtime

[Service]
Type=oneshot
# OpenCode is verified and baked into the immutable ISO. Never try to install
# it into the small live overlay at boot.
ExecStart=/usr/bin/opencode --version
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
